# examples/automated-domain-join/ingresso-lado-cliente-linux/cron.d/gti-ingresso-auto-dominio
#Configuraço para executar o script de ingresso a cada 1 min
* *	* * *	root	/opt/scripts/ingressar-linux-auto.sh
